# .migration-backup/unixsee_campaign_gateway/install/systemd/unixsee-dashboard.service
[Unit]
Description=Unixsee Dashboard
After=network-online.target unixsee-mother.service
Wants=network-online.target

[Service]
Type=simple
User=unixsee
Group=unixsee
WorkingDirectory=/opt/unixsee-campaign-gateway/unixsee_campaign_gateway/dashboard
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=8740
EnvironmentFile=-/etc/unixsee-gateway/dashboard.env
ExecStartPre=/usr/bin/install -d -o unixsee -g unixsee -m 0750 /var/lib/unixsee-gateway/dashboard
ExecStartPre=/usr/bin/install -d -o unixsee -g unixsee -m 0750 /var/log/unixsee-gateway
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 8740
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/unixsee-campaign-gateway/unixsee_campaign_gateway/dashboard/.next/cache /var/lib/unixsee-gateway/dashboard /var/log/unixsee-gateway

[Install]
WantedBy=multi-user.target
